# scripts/com.example.reddit-scout-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.example.reddit-scout-daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/path/to/reddit-scout/scripts/daily_runner.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>33</integer>
    </dict>

    <key>StandardOutPath</key>
    <string><HOME>/reddit-scout-reports/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string><HOME>/reddit-scout-reports/launchd.err.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string><HOME></string>
        <key>USER</key>
        <string><USER></string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
